# dotfiles/config/systemd/user/mail-stephane.service
[Unit]
Description=Retrieves my emails

[Service]
Type=oneshot
ExecStart=/home/yfful/documents/cron/mail.sh
